# chatServer/chat-server.service
[Unit]
Description=Runs Ratchet WebSocket service for Chat website.
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/php /home/collectorsscave/public_html/chatServer/chat-server.php >> /home/collectorsscave/public_html/chatServer/chat-server.log 2>&1 
TimeoutStartSec=0

[Install]
WantedBy=default.target